# certbot-renew.service
[Unit]
Description=Renew Let's Encrypt certs (docker compose) and reload nginx
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/website-host/certbot
ExecStart=/usr/bin/docker compose -f dc-front-proxy.yml run --rm certbot
ExecStartPost=-/usr/bin/docker compose -f dc-front-proxy.yml exec nginx nginx -s reload